# skogmark-go-gen/build/generator.service
[Unit]
Description=generator
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/generator/bin
ExecStart=/bin/sh /opt/generator/bin/bootstrap.sh

User=svip
Group=svip
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target