# skogmark-go-blogger/build/blogger.service
[Unit]
Description=blogger
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/blogger/bin
ExecStart=/bin/sh /opt/blogger/bin/bootstrap.sh

User=svip
Group=svip
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target